# recipes-connectivity/mosquitto/files/mosquitto.cron
* * * * * root cd /root/scripts/mosquitto/ && /usr/bin/python3 publish.py
