# recipes-csgames/config-backup/files/config-backup.service
[Unit]
Description=Configuration Backup Service
After=network.target

[Service]
WorkingDirectory=/opt/config-backup
ExecStart=/usr/bin/python3 /opt/config-backup/app.py
Restart=always
User=config-backup
Group=config-backup

[Install]
WantedBy=multi-user.target
