# recipes-csgames/keystore-admin/files/keystore-admin.cron
* * * * * root cd /root/scripts/keystore-admin/ && /usr/bin/python3 populate.py
